# config/profile.d/java.sh
export JAVA_HOME="/usr/java/jdk1.7.0_04"
export MAVEN_HOME="/opt/apache-maven-3.0.4"
export PLAY_HOME="/opt/play-1.2.4"

export PATH="$JAVA_HOME/bin:$MAVEN_HOME/bin:$PLAY_HOME:$PATH"

